# backend/scripts/liveboost_worker.conf
; ============================================================
; Supervisor config for LiveBoost Analysis Worker
; ============================================================
; Install: sudo apt-get install supervisor
; Copy:    sudo cp liveboost_worker.conf /etc/supervisor/conf.d/
; Reload:  sudo supervisorctl reread && sudo supervisorctl update
; ============================================================

[program:liveboost-worker]
command=python -m app.workers.live_analysis_worker
directory=/home/ubuntu/aitherhub/backend
user=ubuntu
autostart=true
autorestart=true
startretries=10
startsecs=5
stopwaitsecs=600
redirect_stderr=true
stdout_logfile=/var/log/liveboost-worker.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=
    AZURE_STORAGE_CONNECTION_STRING="%(ENV_AZURE_STORAGE_CONNECTION_STRING)s",
    AZURE_QUEUE_NAME="video-jobs",
    DATABASE_URL="%(ENV_DATABASE_URL)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    WORKER_POLL_INTERVAL="5",
    WORKER_MAX_CONCURRENT="2"
